# sv/clipboard-service/run
#!/bin/sh

log_info "Starting clipboard manager (wl-clip-persist)..."

wl-clip-persist --clipboard both
