# sv/xdg-autostart-service/run
#!/bin/bash

log_info "Running 'dex --autostart'"
dex --autostart
