# sv/widget-service/run
#!/bin/bash

eww daemon --no-daemonize &
eww open SystemBar

wait
